# services/seedling-web.service
[Unit]
Description=Seedling web interface
Documentation=https://github.com/beyondessential/seedling
# Order after tailscaled when present so tailscale0 exists at bind time. Ordering
# only: the service still starts if Tailscale isn't installed/enabled (see
# ExecStartPre below).
After=network-online.target seedling.service tailscaled.service
Wants=network-online.target

[Service]
Type=simple
# Runs as root so it can read the daemon-owned bootstrap files under
# /var/lib/seedling (the pinned fingerprint and its own client key). It binds
# only the loopback web port and needs no other privilege.
#
# Authenticates operators via Tailscale identity headers. The HTTP listener
# binds loopback only, so it is reachable exactly via the co-located
# `tailscale serve` front-end (seedling-web-tailscale-serve.service), which
# terminates HTTPS on the tailnet and injects the Tailscale-User-* headers
# seedling-web trusts — nothing untrusted can reach the loopback port to spoof
# them.
#
# The WebTransport listener is bound separately on the tailnet (--wt-interface
# lo,tailscale0): WebTransport is HTTP/3 (QUIC/UDP) and cannot be carried by
# `tailscale serve` (a TCP HTTP proxy), so the browser reaches it directly at
# <node>.<tailnet>.ts.net:7893. It carries its own per-session token and pins
# the server certificate, so exposing it on the tailnet needs no header trust.
# A missing tailscale0 is skipped (not fatal): the interface just isn't served.
#
# For a password-login deployment instead: set auth.password_hash in
# /etc/seedling/web.toml, drop --trust-tailscale-headers with a drop-in, and
# front the HTTP port with a TLS-terminating reverse proxy — but note the
# WebTransport port still must be reachable directly (adjust --wt-interface /
# --wt-listen for that network). See docs/deploying.md.
#
# Best-effort wait so tailscale0 has its address at bind time, WITHOUT blocking
# startup: returns immediately if tailscale isn't installed or tailscaled isn't
# active, waits at most ~30s for the tailnet IP otherwise, and the leading `-`
# means it can never fail the unit. ($$ is a literal $ for the shell.)
ExecStartPre=-/bin/sh -c 'command -v tailscale >/dev/null 2>&1 || exit 0; systemctl is-active --quiet tailscaled 2>/dev/null || exit 0; i=0; while [ "$$i" -lt 30 ]; do tailscale ip -4 >/dev/null 2>&1 && exit 0; i=$$((i+1)); sleep 1; done'
ExecStart=/usr/bin/seedling-web \
    --config /etc/seedling/web.toml \
    --daemon-fingerprint-file /var/lib/seedling/oi.fingerprint \
    --key-file /var/lib/seedling/web.key \
    --trust-tailscale-headers \
    --wt-interface lo,tailscale0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
