# services/seedling-web-tailscale-serve.service
[Unit]
Description=Expose the Seedling web interface over Tailscale (tailscale serve)
Documentation=https://tailscale.com/kb/1242/tailscale-serve
# No-op if Tailscale is not installed: with the condition unmet the unit is
# skipped (not failed) when started, so it is safe to enable regardless.
ConditionPathExists=/usr/bin/tailscale
After=tailscaled.service seedling-web.service
Wants=seedling-web.service

[Service]
# `tailscale serve --bg` persists the mapping in tailscaled and returns, so this
# is a one-shot that stays "active (exited)"; stopping it tears the mapping down.
# tailscaled resumes the mapping across reboots on its own, so a failure here
# (e.g. run before the node is up) does not lose an already-configured serve.
Type=oneshot
RemainAfterExit=yes
# Publish the loopback web port (seedling-web's default 7894) over HTTPS on the
# tailnet. tailscale serve injects the Tailscale-User-* identity headers that
# seedling-web trusts. Served on 7895 (in Seedling's 789x range) — NOT 443,
# which is reserved for app workloads (Caddy).
ExecStart=/usr/bin/tailscale serve --https=7895 --bg localhost:7894
ExecStop=/usr/bin/tailscale serve --https=7895 off

[Install]
WantedBy=multi-user.target
